# php-worker/supervisord.conf
; Supervisor 主設定區塊
[supervisord]
; 設定 Supervisor 以前台模式運行（容器環境必要，避免立即退出）
nodaemon=true
; 設定 Supervisor 本身的日誌路徑
logfile=/var/log/supervisor/supervisord.log
; 設定 Supervisor 日誌大小上限（超過自動輪替）
logfile_maxbytes=50MB
; 設定保留的舊日誌備份份數
logfile_backups=5
; 設定 Supervisor 日誌等級
loglevel=info
; 設定 PID 檔案路徑
pidfile=/tmp/supervisord.pid

; Supervisor RPC 介面設定（supervisorctl 指令需要）
[rpcinterface:supervisor]
; 使用內建的 RPC 工廠
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; supervisorctl 客戶端設定
[supervisorctl]
; 設定 supervisorctl 連線位址（與 unix_http_server 一致）
serverurl=unix:///tmp/supervisor.sock

; Unix socket 伺服器設定（供 supervisorctl 連線使用）
[unix_http_server]
; 設定 Unix socket 檔案路徑
file=/tmp/supervisor.sock
; 設定 socket 檔案權限（允許同群組讀寫）
chmod=0700

; ============================================================
; Worker 程序設定範例（可複製並修改為實際任務）
; ============================================================

; 佇列 Worker 程序設定（以 Laravel Queue Worker 為例）
[program:queue-worker]
; 設定執行的指令（請依專案調整路徑與參數）
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=90
; 設定工作目錄
directory=/var/www/html
; 設定執行使用者（與容器使用者一致）
user=www
; 設定是否自動啟動（Supervisor 啟動時自動開啟此程序）
autostart=true
; 設定程序意外退出後是否自動重啟
autorestart=true
; 設定啟動後多少秒內未退出視為啟動成功
startsecs=3
; 設定最多重試次數（超過後標記為 FATAL）
startretries=3
; 設定程序標準輸出的日誌路徑
stdout_logfile=/var/log/supervisor/queue-worker.log
; 設定標準輸出日誌大小上限
stdout_logfile_maxbytes=10MB
; 設定保留的輸出日誌備份份數
stdout_logfile_backups=3
; 設定程序標準錯誤的日誌路徑
stderr_logfile=/var/log/supervisor/queue-worker-error.log
; 設定標準錯誤日誌大小上限
stderr_logfile_maxbytes=10MB
; 設定保留的錯誤日誌備份份數
stderr_logfile_backups=3

; 排程 Worker 程序設定（以 Laravel Schedule 為例）
[program:schedule-runner]
; 設定每分鐘執行一次排程指令
command=php /var/www/html/artisan schedule:run
; 設定工作目錄
directory=/var/www/html
; 設定執行使用者
user=www
; 設定是否自動啟動
autostart=true
; 設定程序退出後是否自動重啟（排程任務通常執行完即退出）
autorestart=true
; 設定啟動後多少秒視為啟動成功
startsecs=1
; 設定最多重試次數
startretries=3
; 設定標準輸出日誌路徑
stdout_logfile=/var/log/supervisor/schedule-runner.log
; 設定標準輸出日誌大小上限
stdout_logfile_maxbytes=10MB
; 設定保留的輸出日誌備份份數
stdout_logfile_backups=3
; 設定標準錯誤日誌路徑
stderr_logfile=/var/log/supervisor/schedule-runner-error.log
; 設定標準錯誤日誌大小上限
stderr_logfile_maxbytes=10MB
; 設定保留的錯誤日誌備份份數
stderr_logfile_backups=3
